# mf/strongswan.files
usr/lib64/ipsec/libstrongswan.so.0.0.0
usr/lib64/ipsec/libstrongswan.la
usr/lib64/ipsec/libstrongswan.a
usr/lib64/ipsec/plugins/libstrongswan-aes.so
usr/lib64/ipsec/plugins/libstrongswan-aes.la
usr/lib64/ipsec/plugins/libstrongswan-aes.a
usr/lib64/ipsec/plugins/libstrongswan-des.so
usr/lib64/ipsec/plugins/libstrongswan-des.la
usr/lib64/ipsec/plugins/libstrongswan-des.a
usr/lib64/ipsec/plugins/libstrongswan-rc2.so
usr/lib64/ipsec/plugins/libstrongswan-rc2.la
usr/lib64/ipsec/plugins/libstrongswan-rc2.a
usr/lib64/ipsec/plugins/libstrongswan-md5.so
usr/lib64/ipsec/plugins/libstrongswan-md5.la
usr/lib64/ipsec/plugins/libstrongswan-md5.a
usr/lib64/ipsec/plugins/libstrongswan-sha1.so
usr/lib64/ipsec/plugins/libstrongswan-sha1.la
usr/lib64/ipsec/plugins/libstrongswan-sha1.a
usr/lib64/ipsec/plugins/libstrongswan-sha2.so
usr/lib64/ipsec/plugins/libstrongswan-sha2.la
usr/lib64/ipsec/plugins/libstrongswan-sha2.a
usr/lib64/ipsec/plugins/libstrongswan-gmp.so
usr/lib64/ipsec/plugins/libstrongswan-gmp.la
usr/lib64/ipsec/plugins/libstrongswan-gmp.a
usr/lib64/ipsec/plugins/libstrongswan-random.so
usr/lib64/ipsec/plugins/libstrongswan-random.la
usr/lib64/ipsec/plugins/libstrongswan-random.a
usr/lib64/ipsec/plugins/libstrongswan-nonce.so
usr/lib64/ipsec/plugins/libstrongswan-nonce.la
usr/lib64/ipsec/plugins/libstrongswan-nonce.a
usr/lib64/ipsec/plugins/libstrongswan-hmac.so
usr/lib64/ipsec/plugins/libstrongswan-hmac.la
usr/lib64/ipsec/plugins/libstrongswan-hmac.a
usr/lib64/ipsec/plugins/libstrongswan-cmac.so
usr/lib64/ipsec/plugins/libstrongswan-cmac.la
usr/lib64/ipsec/plugins/libstrongswan-cmac.a
usr/lib64/ipsec/plugins/libstrongswan-xcbc.so
usr/lib64/ipsec/plugins/libstrongswan-xcbc.la
usr/lib64/ipsec/plugins/libstrongswan-xcbc.a
usr/lib64/ipsec/plugins/libstrongswan-x509.so
usr/lib64/ipsec/plugins/libstrongswan-x509.la
usr/lib64/ipsec/plugins/libstrongswan-x509.a
usr/lib64/ipsec/plugins/libstrongswan-revocation.so
usr/lib64/ipsec/plugins/libstrongswan-revocation.la
usr/lib64/ipsec/plugins/libstrongswan-revocation.a
usr/lib64/ipsec/plugins/libstrongswan-constraints.so
usr/lib64/ipsec/plugins/libstrongswan-constraints.la
usr/lib64/ipsec/plugins/libstrongswan-constraints.a
usr/lib64/ipsec/plugins/libstrongswan-acert.so
usr/lib64/ipsec/plugins/libstrongswan-acert.la
usr/lib64/ipsec/plugins/libstrongswan-acert.a
usr/lib64/ipsec/plugins/libstrongswan-pubkey.so
usr/lib64/ipsec/plugins/libstrongswan-pubkey.la
usr/lib64/ipsec/plugins/libstrongswan-pubkey.a
usr/lib64/ipsec/plugins/libstrongswan-pkcs1.so
usr/lib64/ipsec/plugins/libstrongswan-pkcs1.la
usr/lib64/ipsec/plugins/libstrongswan-pkcs1.a
usr/lib64/ipsec/plugins/libstrongswan-pkcs7.so
usr/lib64/ipsec/plugins/libstrongswan-pkcs7.la
usr/lib64/ipsec/plugins/libstrongswan-pkcs7.a
usr/lib64/ipsec/plugins/libstrongswan-pkcs8.so
usr/lib64/ipsec/plugins/libstrongswan-pkcs8.la
usr/lib64/ipsec/plugins/libstrongswan-pkcs8.a
usr/lib64/ipsec/plugins/libstrongswan-pkcs12.so
usr/lib64/ipsec/plugins/libstrongswan-pkcs12.la
usr/lib64/ipsec/plugins/libstrongswan-pkcs12.a
usr/lib64/ipsec/plugins/libstrongswan-pgp.so
usr/lib64/ipsec/plugins/libstrongswan-pgp.la
usr/lib64/ipsec/plugins/libstrongswan-pgp.a
usr/lib64/ipsec/plugins/libstrongswan-dnskey.so
usr/lib64/ipsec/plugins/libstrongswan-dnskey.la
usr/lib64/ipsec/plugins/libstrongswan-dnskey.a
usr/lib64/ipsec/plugins/libstrongswan-sshkey.so
usr/lib64/ipsec/plugins/libstrongswan-sshkey.la
usr/lib64/ipsec/plugins/libstrongswan-sshkey.a
usr/lib64/ipsec/plugins/libstrongswan-pem.so
usr/lib64/ipsec/plugins/libstrongswan-pem.la
usr/lib64/ipsec/plugins/libstrongswan-pem.a
usr/lib64/ipsec/plugins/libstrongswan-curl.so
usr/lib64/ipsec/plugins/libstrongswan-curl.la
usr/lib64/ipsec/plugins/libstrongswan-curl.a
usr/lib64/ipsec/plugins/libstrongswan-files.so
usr/lib64/ipsec/plugins/libstrongswan-files.la
usr/lib64/ipsec/plugins/libstrongswan-files.a
usr/lib64/ipsec/plugins/libstrongswan-openssl.so
usr/lib64/ipsec/plugins/libstrongswan-openssl.la
usr/lib64/ipsec/plugins/libstrongswan-openssl.a
usr/lib64/ipsec/plugins/libstrongswan-fips-prf.so
usr/lib64/ipsec/plugins/libstrongswan-fips-prf.la
usr/lib64/ipsec/plugins/libstrongswan-fips-prf.a
usr/lib64/ipsec/plugins/libstrongswan-pkcs11.so
usr/lib64/ipsec/plugins/libstrongswan-pkcs11.la
usr/lib64/ipsec/plugins/libstrongswan-pkcs11.a
usr/lib64/ipsec/plugins/libstrongswan-socket-default.so
usr/lib64/ipsec/plugins/libstrongswan-socket-default.la
usr/lib64/ipsec/plugins/libstrongswan-socket-default.a
usr/lib64/ipsec/plugins/libstrongswan-connmark.so
usr/lib64/ipsec/plugins/libstrongswan-connmark.la
usr/lib64/ipsec/plugins/libstrongswan-connmark.a
usr/lib64/ipsec/plugins/libstrongswan-forecast.so
usr/lib64/ipsec/plugins/libstrongswan-forecast.la
usr/lib64/ipsec/plugins/libstrongswan-forecast.a
usr/lib64/ipsec/plugins/libstrongswan-farp.so
usr/lib64/ipsec/plugins/libstrongswan-farp.la
usr/lib64/ipsec/plugins/libstrongswan-farp.a
usr/lib64/ipsec/plugins/libstrongswan-stroke.so
usr/lib64/ipsec/plugins/libstrongswan-stroke.la
usr/lib64/ipsec/plugins/libstrongswan-stroke.a
usr/lib64/ipsec/plugins/libstrongswan-vici.so
usr/lib64/ipsec/plugins/libstrongswan-vici.la
usr/lib64/ipsec/plugins/libstrongswan-vici.a
usr/lib64/ipsec/plugins/libstrongswan-updown.so
usr/lib64/ipsec/plugins/libstrongswan-updown.la
usr/lib64/ipsec/plugins/libstrongswan-updown.a
usr/lib64/ipsec/plugins/libstrongswan-ext-auth.so
usr/lib64/ipsec/plugins/libstrongswan-ext-auth.la
usr/lib64/ipsec/plugins/libstrongswan-ext-auth.a
usr/lib64/ipsec/plugins/libstrongswan-eap-sim.so
usr/lib64/ipsec/plugins/libstrongswan-eap-sim.la
usr/lib64/ipsec/plugins/libstrongswan-eap-sim.a
usr/lib64/ipsec/plugins/libstrongswan-eap-sim-file.so
usr/lib64/ipsec/plugins/libstrongswan-eap-sim-file.la
usr/lib64/ipsec/plugins/libstrongswan-eap-sim-file.a
usr/lib64/ipsec/plugins/libstrongswan-eap-sim-pcsc.so
usr/lib64/ipsec/plugins/libstrongswan-eap-sim-pcsc.la
usr/lib64/ipsec/plugins/libstrongswan-eap-sim-pcsc.a
usr/lib64/ipsec/plugins/libstrongswan-eap-tls.so
usr/lib64/ipsec/plugins/libstrongswan-eap-tls.la
usr/lib64/ipsec/plugins/libstrongswan-eap-tls.a
usr/lib64/ipsec/plugins/libstrongswan-eap-ttls.so
usr/lib64/ipsec/plugins/libstrongswan-eap-ttls.la
usr/lib64/ipsec/plugins/libstrongswan-eap-ttls.a
usr/lib64/ipsec/plugins/libstrongswan-dhcp.so
usr/lib64/ipsec/plugins/libstrongswan-dhcp.la
usr/lib64/ipsec/plugins/libstrongswan-dhcp.a
usr/lib64/ipsec/plugins/libstrongswan-kernel-netlink.so
usr/lib64/ipsec/plugins/libstrongswan-kernel-netlink.la
usr/lib64/ipsec/plugins/libstrongswan-kernel-netlink.a
usr/lib64/ipsec/plugins/libstrongswan-kernel-libipsec.so
usr/lib64/ipsec/plugins/libstrongswan-kernel-libipsec.la
usr/lib64/ipsec/plugins/libstrongswan-kernel-libipsec.a
usr/lib64/ipsec/plugins/libstrongswan-whitelist.so
usr/lib64/ipsec/plugins/libstrongswan-whitelist.la
usr/lib64/ipsec/plugins/libstrongswan-whitelist.a
usr/lib64/ipsec/plugins/libstrongswan-certexpire.so
usr/lib64/ipsec/plugins/libstrongswan-certexpire.la
usr/lib64/ipsec/plugins/libstrongswan-certexpire.a
usr/lib64/ipsec/plugins/libstrongswan-duplicheck.so
usr/lib64/ipsec/plugins/libstrongswan-duplicheck.la
usr/lib64/ipsec/plugins/libstrongswan-duplicheck.a
usr/lib64/ipsec/plugins/libstrongswan-xauth-generic.so
usr/lib64/ipsec/plugins/libstrongswan-xauth-generic.la
usr/lib64/ipsec/plugins/libstrongswan-xauth-generic.a
usr/lib64/ipsec/plugins/libstrongswan-xauth-eap.so
usr/lib64/ipsec/plugins/libstrongswan-xauth-eap.la
usr/lib64/ipsec/plugins/libstrongswan-xauth-eap.a
usr/lib64/ipsec/plugins/libstrongswan-xauth-pam.so
usr/lib64/ipsec/plugins/libstrongswan-xauth-pam.la
usr/lib64/ipsec/plugins/libstrongswan-xauth-pam.a
usr/lib64/ipsec/plugins/libstrongswan-xauth-noauth.so
usr/lib64/ipsec/plugins/libstrongswan-xauth-noauth.la
usr/lib64/ipsec/plugins/libstrongswan-xauth-noauth.a
usr/lib64/ipsec/plugins/libstrongswan-resolve.so
usr/lib64/ipsec/plugins/libstrongswan-resolve.la
usr/lib64/ipsec/plugins/libstrongswan-resolve.a
usr/lib64/ipsec/plugins/libstrongswan-attr.so
usr/lib64/ipsec/plugins/libstrongswan-attr.la
usr/lib64/ipsec/plugins/libstrongswan-attr.a
usr/lib64/ipsec/libipsec.so.0.0.0
usr/lib64/ipsec/libipsec.la
usr/lib64/ipsec/libipsec.a
usr/lib64/ipsec/libsimaka.so.0.0.0
usr/lib64/ipsec/libsimaka.la
usr/lib64/ipsec/libsimaka.a
usr/lib64/ipsec/libtls.so.0.0.0
usr/lib64/ipsec/libtls.la
usr/lib64/ipsec/libtls.a
usr/lib64/ipsec/libcharon.so.0.0.0
usr/lib64/ipsec/libcharon.la
usr/lib64/ipsec/libcharon.a
usr/lib64/ipsec/libvici.so.0.0.0
usr/lib64/ipsec/libvici.la
usr/lib64/ipsec/libvici.a
usr/libexec/ipsec/whitelist
usr/libexec/ipsec/duplicheck
usr/libexec/ipsec/starter
usr/libexec/ipsec/_copyright
usr/libexec/ipsec/charon
usr/libexec/ipsec/charon-nm
usr/libexec/ipsec/stroke
usr/libexec/ipsec/_updown
usr/libexec/ipsec/scepclient
usr/sbin/ipsec
usr/sbin/charon-systemd
usr/sbin/swanctl
usr/share/man/man8/ipsec.8
usr/share/man/man8/scepclient.8
usr/share/man/man8/swanctl.8
usr/share/man/man1/pki.1
usr/share/man/man1/pki---acert.1
usr/share/man/man1/pki---dn.1
usr/share/man/man1/pki---gen.1
usr/share/man/man1/pki---issue.1
usr/share/man/man1/pki---keyid.1
usr/share/man/man1/pki---pkcs7.1
usr/share/man/man1/pki---print.1
usr/share/man/man1/pki---pub.1
usr/share/man/man1/pki---req.1
usr/share/man/man1/pki---self.1
usr/share/man/man1/pki---signcrl.1
usr/share/man/man1/pki---verify.1
usr/share/man/man5/swanctl.conf.5
usr/share/man/man5/ipsec.conf.5
usr/share/man/man5/ipsec.secrets.5
usr/share/man/man5/strongswan.conf.5
usr/share/strongswan/templates/config/strongswan.d/charon.conf
usr/share/strongswan/templates/config/strongswan.d/charon-logging.conf
usr/share/strongswan/templates/config/strongswan.d/starter.conf
usr/share/strongswan/templates/config/strongswan.d/scepclient.conf
usr/share/strongswan/templates/config/strongswan.d/pki.conf
usr/share/strongswan/templates/config/strongswan.d/swanctl.conf
usr/share/strongswan/templates/config/plugins/pkcs11.conf
usr/share/strongswan/templates/config/plugins/aes.conf
usr/share/strongswan/templates/config/plugins/des.conf
usr/share/strongswan/templates/config/plugins/rc2.conf
usr/share/strongswan/templates/config/plugins/sha2.conf
usr/share/strongswan/templates/config/plugins/sha1.conf
usr/share/strongswan/templates/config/plugins/md5.conf
usr/share/strongswan/templates/config/plugins/random.conf
usr/share/strongswan/templates/config/plugins/nonce.conf
usr/share/strongswan/templates/config/plugins/x509.conf
usr/share/strongswan/templates/config/plugins/revocation.conf
usr/share/strongswan/templates/config/plugins/constraints.conf
usr/share/strongswan/templates/config/plugins/acert.conf
usr/share/strongswan/templates/config/plugins/pubkey.conf
usr/share/strongswan/templates/config/plugins/pkcs1.conf
usr/share/strongswan/templates/config/plugins/pkcs7.conf
usr/share/strongswan/templates/config/plugins/pkcs8.conf
usr/share/strongswan/templates/config/plugins/pkcs12.conf
usr/share/strongswan/templates/config/plugins/pgp.conf
usr/share/strongswan/templates/config/plugins/dnskey.conf
usr/share/strongswan/templates/config/plugins/sshkey.conf
usr/share/strongswan/templates/config/plugins/pem.conf
usr/share/strongswan/templates/config/plugins/openssl.conf
usr/share/strongswan/templates/config/plugins/fips-prf.conf
usr/share/strongswan/templates/config/plugins/gmp.conf
usr/share/strongswan/templates/config/plugins/xcbc.conf
usr/share/strongswan/templates/config/plugins/cmac.conf
usr/share/strongswan/templates/config/plugins/hmac.conf
usr/share/strongswan/templates/config/plugins/curl.conf
usr/share/strongswan/templates/config/plugins/files.conf
usr/share/strongswan/templates/config/plugins/attr.conf
usr/share/strongswan/templates/config/plugins/kernel-libipsec.conf
usr/share/strongswan/templates/config/plugins/kernel-netlink.conf
usr/share/strongswan/templates/config/plugins/resolve.conf
usr/share/strongswan/templates/config/plugins/socket-default.conf
usr/share/strongswan/templates/config/plugins/connmark.conf
usr/share/strongswan/templates/config/plugins/forecast.conf
usr/share/strongswan/templates/config/plugins/farp.conf
usr/share/strongswan/templates/config/plugins/stroke.conf
usr/share/strongswan/templates/config/plugins/vici.conf
usr/share/strongswan/templates/config/plugins/updown.conf
usr/share/strongswan/templates/config/plugins/eap-sim.conf
usr/share/strongswan/templates/config/plugins/eap-sim-file.conf
usr/share/strongswan/templates/config/plugins/eap-sim-pcsc.conf
usr/share/strongswan/templates/config/plugins/eap-tls.conf
usr/share/strongswan/templates/config/plugins/eap-ttls.conf
usr/share/strongswan/templates/config/plugins/xauth-generic.conf
usr/share/strongswan/templates/config/plugins/xauth-eap.conf
usr/share/strongswan/templates/config/plugins/xauth-pam.conf
usr/share/strongswan/templates/config/plugins/xauth-noauth.conf
usr/share/strongswan/templates/config/plugins/dhcp.conf
usr/share/strongswan/templates/config/plugins/whitelist.conf
usr/share/strongswan/templates/config/plugins/ext-auth.conf
usr/share/strongswan/templates/config/plugins/certexpire.conf
usr/share/strongswan/templates/config/plugins/duplicheck.conf
usr/share/strongswan/templates/config/strongswan.conf
usr/bin/pki
usr/lib/systemd/system/strongswan.service
usr/lib/systemd/system/strongswan-swanctl.service
opt/strongswan/etc/ipsec.conf
opt/strongswan/etc/ipsec.secrets
opt/strongswan/etc/dbus-1/system.d/nm-strongswan-service.conf
opt/strongswan/etc/swanctl/swanctl.conf
opt/strongswan/etc/strongswan.d/charon/pkcs11.conf
opt/strongswan/etc/strongswan.d/charon/aes.conf
opt/strongswan/etc/strongswan.d/charon/des.conf
opt/strongswan/etc/strongswan.d/charon/rc2.conf
opt/strongswan/etc/strongswan.d/charon/sha2.conf
opt/strongswan/etc/strongswan.d/charon/sha1.conf
opt/strongswan/etc/strongswan.d/charon/md5.conf
opt/strongswan/etc/strongswan.d/charon/random.conf
opt/strongswan/etc/strongswan.d/charon/nonce.conf
opt/strongswan/etc/strongswan.d/charon/x509.conf
opt/strongswan/etc/strongswan.d/charon/revocation.conf
opt/strongswan/etc/strongswan.d/charon/constraints.conf
opt/strongswan/etc/strongswan.d/charon/acert.conf
opt/strongswan/etc/strongswan.d/charon/pubkey.conf
opt/strongswan/etc/strongswan.d/charon/pkcs1.conf
opt/strongswan/etc/strongswan.d/charon/pkcs7.conf
opt/strongswan/etc/strongswan.d/charon/pkcs8.conf
opt/strongswan/etc/strongswan.d/charon/pkcs12.conf
opt/strongswan/etc/strongswan.d/charon/pgp.conf
opt/strongswan/etc/strongswan.d/charon/dnskey.conf
opt/strongswan/etc/strongswan.d/charon/sshkey.conf
opt/strongswan/etc/strongswan.d/charon/pem.conf
opt/strongswan/etc/strongswan.d/charon/openssl.conf
opt/strongswan/etc/strongswan.d/charon/fips-prf.conf
opt/strongswan/etc/strongswan.d/charon/gmp.conf
opt/strongswan/etc/strongswan.d/charon/xcbc.conf
opt/strongswan/etc/strongswan.d/charon/cmac.conf
opt/strongswan/etc/strongswan.d/charon/hmac.conf
opt/strongswan/etc/strongswan.d/charon/curl.conf
opt/strongswan/etc/strongswan.d/charon/files.conf
opt/strongswan/etc/strongswan.d/charon/attr.conf
opt/strongswan/etc/strongswan.d/charon/kernel-libipsec.conf
opt/strongswan/etc/strongswan.d/charon/kernel-netlink.conf
opt/strongswan/etc/strongswan.d/charon/resolve.conf
opt/strongswan/etc/strongswan.d/charon/socket-default.conf
opt/strongswan/etc/strongswan.d/charon/connmark.conf
opt/strongswan/etc/strongswan.d/charon/forecast.conf
opt/strongswan/etc/strongswan.d/charon/farp.conf
opt/strongswan/etc/strongswan.d/charon/stroke.conf
opt/strongswan/etc/strongswan.d/charon/vici.conf
opt/strongswan/etc/strongswan.d/charon/updown.conf
opt/strongswan/etc/strongswan.d/charon/eap-sim.conf
opt/strongswan/etc/strongswan.d/charon/eap-sim-file.conf
opt/strongswan/etc/strongswan.d/charon/eap-sim-pcsc.conf
opt/strongswan/etc/strongswan.d/charon/eap-tls.conf
opt/strongswan/etc/strongswan.d/charon/eap-ttls.conf
opt/strongswan/etc/strongswan.d/charon/xauth-generic.conf
opt/strongswan/etc/strongswan.d/charon/xauth-eap.conf
opt/strongswan/etc/strongswan.d/charon/xauth-pam.conf
opt/strongswan/etc/strongswan.d/charon/xauth-noauth.conf
opt/strongswan/etc/strongswan.d/charon/dhcp.conf
opt/strongswan/etc/strongswan.d/charon/whitelist.conf
opt/strongswan/etc/strongswan.d/charon/ext-auth.conf
opt/strongswan/etc/strongswan.d/charon/certexpire.conf
opt/strongswan/etc/strongswan.d/charon/duplicheck.conf
opt/strongswan/etc/strongswan.d/charon.conf
opt/strongswan/etc/strongswan.d/charon-logging.conf
opt/strongswan/etc/strongswan.d/starter.conf
opt/strongswan/etc/strongswan.d/scepclient.conf
opt/strongswan/etc/strongswan.d/pki.conf
opt/strongswan/etc/strongswan.d/swanctl.conf
opt/strongswan/etc/strongswan.conf
